# docker/rootfs/etc/services.d/nginx/run
#!/bin/bash
LD_LIBRARY_PATH=/usr/local/nginx/lib /usr/local/nginx/sbin/nginx